# piwebcam.service
[Unit]
Description=PiWebcam Streaming Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/PiWebcam
ExecStart=/usr/bin/python3 /home/pi/PiWebcam/webcam.py --quality 50
# Customize arguments above:
#   --quality 50              # JPEG quality 1-100 (lower=faster, default: 85)
#   --resolution 640x480      # Camera resolution (default: 640x480)
#   --framerate 30            # Camera framerate (default: 30)
#   --port 8000               # HTTP port (default: 8000)
#   --ssl                     # Enable HTTPS (requires cert.pem and key.pem)
#   --cert cert.pem           # SSL certificate path
#   --key key.pem             # SSL key path
#   --motion-detect           # Enable motion detection
#   --motion-snapshot         # Save snapshots on motion
#   --motion-snapshot-limit 50  # Max snapshots in RAM
# Example with HTTPS: ExecStart=/usr/bin/python3 /home/pi/PiWebcam/webcam.py --ssl --port 8443 --quality 50
# Example with motion: ExecStart=/usr/bin/python3 /home/pi/PiWebcam/webcam.py --quality 50 --motion-detect --motion-snapshot --motion-snapshot-limit 100
Restart=on-failure
RestartSec=5

# Optional: Enable authentication via environment variables
# Environment="WEBCAM_USER=admin"
# Environment="WEBCAM_PASS=your_password"

# Or load from environment file (recommended for credentials)
# EnvironmentFile=/etc/piwebcam/credentials.env

# Security hardening
# NoNewPrivileges=true
# PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=piwebcam

[Install]
WantedBy=multi-user.target
